# deploy/systemd/tradepulse.service
[Unit]
Description=TradePulse FastAPI application server
Documentation=https://github.com/TradePulse/TradePulse/blob/main/DEPLOYMENT.md
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=tradepulse
Group=tradepulse
WorkingDirectory=/opt/tradepulse/current
EnvironmentFile=-/etc/tradepulse/tradepulse.env
Environment="APP_HOME=/opt/tradepulse/current"
Environment="VENV_DIR=/opt/tradepulse/venv"
ExecStart=/usr/local/bin/start-tradepulse.sh
Restart=always
RestartSec=5s
KillMode=mixed
TimeoutStopSec=45s
LimitNOFILE=65535
SyslogIdentifier=tradepulse-api

[Install]
WantedBy=multi-user.target
